# deploy/systemd/docomator-api.service.in
[Unit]
Description=Docomator API
After=network.target
Wants=network.target

[Service]
Type=simple
User=@DOCOMATOR_USER@
Group=@DOCOMATOR_GROUP@
EnvironmentFile=@DOCOMATOR_CONFIG_DIR@/docomator.env
Environment=DOCOMATOR_RELEASE_METADATA_PATH=@DOCOMATOR_INSTALL_ROOT@/current/release.json
WorkingDirectory=@DOCOMATOR_INSTALL_ROOT@/current/app
ExecStartPre=@DOCOMATOR_INSTALL_ROOT@/current/runtime/node/bin/node @DOCOMATOR_INSTALL_ROOT@/current/app/scripts/runtime/migrate.mjs
ExecStart=@DOCOMATOR_INSTALL_ROOT@/current/runtime/node/bin/node @DOCOMATOR_INSTALL_ROOT@/current/app/apps/api/dist/server.js
Restart=on-failure
RestartSec=3
TimeoutStopSec=15
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
ReadWritePaths=@DOCOMATOR_DATA_DIR@

[Install]
WantedBy=multi-user.target
